# cookbooks/base-common/files/default/etc/cron.d/chef-client
# Run chef-client every 5 minutes
*/5 * * * * root chef-client

